# configs/supervisor/ser.conf
[program:logfile_consumer]
command=/wwwroot/amustware/ser/app/console rabbitmq:consumer -w -l 256 -m 100 upload_log
user=www-data
autorestart=true
autostart=true